# deploy/systemd/mbo-backfill.timer
[Unit]
Description=Run mbo-backfill nightly at 04:00 (catches up missed nights on boot)

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
